# debian/prerm
#!/bin/sh
set -e
case "$1" in
	upgrade|remove)
		invoke-rc.d etcd stop
	;;
esac
exit 0
